# init.d/smwd
#!/bin/sh /etc/rc.common
# smwd init.d script

START=80
STOP=15

start() {
	/usr/local/smwd.sh &
}

stop() {
	kill $(pidof smwd.sh)
}
