# init.d/smwd
#!/bin/sh /etc/rc.common
# mesh-wd init.d script

START=20
STOP=15

start() {
	/usr/local/smwd.sh &
}

stop() {
	kill $(pidof smwd.sh)
}
